# systemd/attestd.service
[Unit]
Description=CorridorOS Attestation Daemon
Documentation=https://github.com/corridoros/corridoros
After=network.target
Wants=network.target

[Service]
Type=simple
User=corridoros
Group=corridoros
ExecStart=/usr/local/bin/attestd
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=attestd

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/corridoros
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment=ATTESTD_LOG_LEVEL=info
Environment=ATTESTD_CONFIG_PATH=/etc/corridoros/attestd.conf

[Install]
WantedBy=multi-user.target
