# systemd/corrd.service
[Unit]
Description=CorridorOS Photonic Corridor Management Daemon
Documentation=https://github.com/corridoros/corridoros
After=network.target
Wants=network.target

[Service]
Type=simple
User=corridoros
Group=corridoros
ExecStart=/usr/local/bin/corrd
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=corrd

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/corridoros
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment=CORRD_LOG_LEVEL=info
Environment=CORRD_CONFIG_PATH=/etc/corridoros/corrd.conf

[Install]
WantedBy=multi-user.target
